# myearth.service
[Unit]
Description=MyEarth FastAPI Application
After=network.target

[Service]
Type=simple
User=jc
WorkingDirectory=/home/jc/MyEarth
Environment=PATH=/home/jc/MyEarth/venv/bin
Environment=PORT=5000
ExecStart=/home/jc/MyEarth/venv/bin/python main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 