# infrastructure/myearth.service
[Unit]
Description=MyEarth FastAPI Application
After=network.target

[Service]
Type=simple
User=jc
WorkingDirectory=/home/jc/MyEarth
Environment=PATH=/home/jc/MyEarth/venv/bin
Environment=PORT=5001
Environment=PYTHONPATH=/home/jc/MyEarth
Environment=DATABASE_URL=postgresql://postgres:myearth_password@localhost/myearth
ExecStart=/home/jc/MyEarth/venv/bin/python main.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=10
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target 